# _superninja_startup.conf
[program:7860_python3.11]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec /usr/local/bin/python3.11 /usr/local/bin/gunicorn --bind 0.0.0.0:7860 --workers 2 --threads 4 --timeout 300 --keep-alive 5 --log-level info app:app"
directory=/workspace/AI_grant_writer
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/7860_python3.11.err.log
stdout_logfile=/var/log/supervisor/7860_python3.11.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
